# supervisord.conf
[supervisord]
nodaemon=true
logfile=/mnt/data/openguardrails-data/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/mnt/data/openguardrails-data/logs

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/mnt/data/openguardrails-data/logs/nginx.log
stderr_logfile=/mnt/data/openguardrails-data/logs/nginx_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:admin-service]
command=/app/supervisor-entrypoint.sh python3 start_admin_service.py
directory=/app
autostart=true
autorestart=true
priority=20
environment=SERVICE_NAME="admin"
stdout_logfile=/mnt/data/openguardrails-data/logs/admin.log
stderr_logfile=/mnt/data/openguardrails-data/logs/admin_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:detection-service]
command=/app/supervisor-entrypoint.sh python3 start_detection_service.py
directory=/app
autostart=true
autorestart=true
priority=30
environment=SERVICE_NAME="detection"
stdout_logfile=/mnt/data/openguardrails-data/logs/detection.log
stderr_logfile=/mnt/data/openguardrails-data/logs/detection_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:proxy-service]
command=/app/supervisor-entrypoint.sh python3 start_proxy_service.py
directory=/app
autostart=true
autorestart=true
priority=40
environment=SERVICE_NAME="proxy"
stdout_logfile=/mnt/data/openguardrails-data/logs/proxy.log
stderr_logfile=/mnt/data/openguardrails-data/logs/proxy_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
